# recipes-extended/cyberbot-bringup/files/cyberbot.service
[Unit]
Description=CyberBot Service

[Service]
ExecStart=/usr/bin/cyberbot-init
Restart=always
User=root

[Install]
WantedBy=multi-user.target